# 配置/supervisor/ini/ally.ini
[program:ally]
command=/work/python37/bin/gunicorn ally_backend.wsgi:application -b 0.0.0.0:8000

directory=/work/ally/ally_backend/
startsecs=1
stopwaitsecs=1
autostart=false
autorestart=false
user=root
stdout_logfile=/work/log/ally/ally.log
stderr_logfile=/work/log/ally/ally.err
