# 配置/supervisor/ini/celery-beat.ini
[program:celery-beat]
command=celery -A ally_backend beat -l debug

directory=/work/ally/ally_backend/
startsecs=1
stopwaitsecs=600
autostart=false
autorestart=false
user=root
stdout_logfile=/work/log/ally/celery-beat.log
stderr_logfile=/work/log/ally/celery-beat.err
