# 配置/supervisor/ini/celery-fk.ini
[program:celery_fk]
command=celery -A fk_backend worker -l debug

directory=/work/floopyknight/server/
startsecs=1
stopwaitsecs=600
autostart=false
autorestart=false
user=root
stdout_logfile=/work/log/fk_log/celery.log
stderr_logfile=/work/log/fk_log/celery.err
