# 配置/supervisor/ini/celery.ini
[program:celery]
command=celery -A ally_backend worker -l debug

directory=/work/ally/ally_backend/
startsecs=1
stopwaitsecs=600
autostart=false
autorestart=false
user=root
stdout_logfile=/work/log/ally/celery.log
stderr_logfile=/work/log/ally/celery.err
