# 配置/supervisor/ini/floopyknight.ini
[program:fk]
command=/work/python37/bin/gunicorn fk_backend.wsgi:application -b 0.0.0.0:8002

directory=/work/floopyknight/server/
startsecs=1
stopwaitsecs=1
autostart=false
autorestart=false
user=root
stdout_logfile=/work/log/fk_log/fk.log
stderr_logfile=/work/log/fk_log/fk.err
