# 配置/supervisor/ini/yb.ini
[program:yb]
command=/work/python37/bin/gunicorn yb_backend.wsgi:application -b 0.0.0.0:8001

directory=/work/youbao/yb_backend/
startsecs=1
stopwaitsecs=1
autostart=false
autorestart=false
user=root
stdout_logfile=/work/log/yb/yb.log
stderr_logfile=/work/log/yb/yb.err
